# rpi/tailscale/tailscale-update.service
[Unit]
Description=Run tailscale update

[Service]
Type=oneshot
ExecStart=tailscale update
