# rpi/vault-agent/pihole-restart.path
[Path]
PathChanged=/home/pi/work/vault/certs/pihole.pem

[Install]
WantedBy=multi-user.target
