# rpi/vault-agent/pihole-restart.service
[Service]
Type=oneshot
ExecStart=service lighttpd restart
